# supervisor.d/service.ini
[program:crond]
command         = /usr/sbin/crond -f
directory       = /app
autostart       = true
autorestart     = true
startsecs       = 5
numprocs        = 1
stopsignal      = INT
stdout_logfile  = /var/log/supervisord.log
redirect_stderr = true

[program:bird]
command         = /usr/sbin/bird -c /app/bird.conf -f 
directory       = /app
autostart       = true
autorestart     = true
startsecs       = 5
numprocs        = 1
stopsignal      = INT
stdout_logfile  = /var/log/supervisord.log
redirect_stderr = true
